# rootfs/etc/s6/cron/run
#!/bin/bash

if [ -f ./setup ]
then
  ./setup
fi

/usr/sbin/crond -f
