# rootfs/etc/s6/cron/run
#!/bin/bash

if [ -f ./setup ]
then
  source ./setup
fi

/usr/sbin/crond -f -l 2
