# archive/v3.14/overlay/etc/s6/cron/run
#!/bin/bash

if [ -f ./setup ]; then
    source ./setup
fi

echo "Starting cron daemon now..."
/usr/sbin/crond -f -l 2
